# testnet/ansible/roles/signetd/templates/signetd.service.j2
[Unit]
Description=Signet threshold signing node
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ signetd_user }}
Group={{ signetd_user }}
ExecStart={{ signetd_bin }} -config {{ signetd_config_dir }}/config.yaml -log-level {{ log_level }}
Restart=on-failure
RestartSec=5
LimitNOFILE=65536

# Logging to journald (also available via journalctl -u signetd)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=signetd

[Install]
WantedBy=multi-user.target
